# deploy/primehunt-worker@.service
[Unit]
Description=Primehunt Worker (%i)
After=network.target

[Service]
Type=simple
# Override ExecStart via drop-in or edit this template with actual search parameters.
# Example: primehunt --coordinator http://COORD:8080 --worker-id %H kbn --k 3 --base 2 --min-n 1000 --max-n 100000
EnvironmentFile=/opt/primehunt/.env
ExecStart=/usr/local/bin/primehunt --coordinator http://COORDINATOR:8080 --worker-id %H-%i --checkpoint /opt/primehunt/primehunt-%i.checkpoint kbn --k 3 --base 2 --min-n 1000 --max-n 100000
WorkingDirectory=/opt/primehunt
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
